# contrib/systemd/hive-reap@.timer
[Unit]
Description=Hive reap timer: %I
PartOf=hive@%i.target

[Timer]
OnCalendar=hourly
Persistent=true
